# pushipv6.timer
[Timer]
OnCalendar=*-*-* 08:00:00
Unit=pushipv6.service

[Unit]
Description=IPv6 push Timer 

[Install]
WantedBy=multi-user.target